# bitrix/crontab/crontab.cfg
0,30 * * * * /usr/bin/php -f /home/bitrix/www/local/php_interface/cron_search_reindex.php
*/3 * * * * /usr/bin/php -f /home/bitrix/www/mibix_yandex_cron.php

* * * * * php -f /home/bitrix/www/bitrix/modules/acrit.exportpro/tools/cron_events.php
